# devops_practical_repo/docker-onprem/nginx-compose.service
[Unit]
Description=Docker Compose NGINX App
