# devops_practical_assignment/nginx-compose.service
[Unit]
Description=Docker Compose NGINX App
After=network.target docker.service
Requires=docker.service

[Service]
WorkingDirectory=/opt/nginx-app
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down
Restart=always

[Install]
WantedBy=multi-user.target
